# nbjardins/templates/nbjardins.service.j2
[Unit]
Description=nbjardins
After=syslog.target

[Service]
User=root
ExecStart={{ APPLICATION_DIR }}/spring-app.jar
SuccessExitStatus=143


[Install]
WantedBy=multi-user.target
